# docker/conf/spark-worker/supervisord/supervisord.conf
[supervisord]
nodaemon=true

[program:spark]
command=/usr/local/spark/sbin/start-slave.sh spark://spark-master:7077
