# docker/conf/hadoop-slave/supervisord/supervisord.conf
[supervisord]
nodaemon=true

[program:hadoop]
command=echo "hadoop"
