# docker/conf/hadoop-master/supervisord/supervisord.conf
[supervisord]
nodaemon=true

[program:ssh]
priority=100
autorestart=false
command=service ssh start

[program:hadoop]
priority=200
command=start-dfs.sh
