# docker/conf/probe/supervisord/supervisord.conf
[supervisord]
nodaemon=true

[program:ssh]
autorestart=false
command=service ssh start