# docker/conf/spark-master/supervisord/supervisord.conf
[supervisord]
nodaemon=true

[program:spark]
command=/usr/local/spark/sbin/start-master.sh
